# modules/restic/templates/restic-file.service
[Unit]
Description=restic file [<%= @title %>]

[Service]
Type=oneshot
#Restart=always
EnvironmentFile=/etc/restic/env
ExecStart=restic backup --tag <%= @backup_tag %> <%= @extra_flags %> "<%= @directory %>"

