# modules/restic/templates/restic-mariadb.timer
[Unit]
Description=Run restic mariadb

[Timer]
#OnBootSec=5min
#OnUnitActiveSec=5m
RandomizedDelaySec=<%= @randomized_delay %>
# run if host was restarted in meantime
Persistent=true
#systemd-analyze calendar --iterations=5 *-*-* 3:00 to validate
OnCalendar=<%= @backup_schedule %>
Unit=restic-mariadb.service

[Install]
WantedBy=timers.target