# modules/restic/templates/restic-maintenance.service
[Unit]
Description=restic backup maintenance

[Service]
Type=simple
#Restart=always
ExecStart=/etc/restic/maintenance.sh
