# modules/minio/templates/minio.service
# <%=  __FILE__.gsub(/.*?modules\//,'puppet://modules/') %>

[Unit]
Description=Minio S3 server
ConditionPathExists=<%= @data_dir %>


[Service]
Type=simple
User=minio
Group=minio
Environment=MINIO_ROOT_USER=<%= @access_key %>
Environment=MINIO_ROOT_PASSWORD=<%= @secret_key %>
ExecStart=/usr/local/bin/minio server "<%= @data_dir %>"
Restart=always
RestartSec=30
LimitNOFILE=40000

[Install]
WantedBy=multi-user.target

